# moidify.service
[Unit]
Description=Moidify — Self-hosted music server
After=network.target

[Service]
Type=simple
User=moidify
Group=moidify
WorkingDirectory=/opt/moidify
ExecStart=/opt/moidify/venv/bin/uvicorn server:app --host 0.0.0.0 --port 8000
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5
TimeoutStopSec=30
KillMode=mixed
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
